# OS/lab8/part1/timeSaver.service
[Unit]
Description=timeSaver

[Service]
Type=simple
ExecStart=/home/fnick/time

[Install]
WantedBy=multi-user.target
